# powershell/Invoke-LogParser/tests/samples/defender-alerts.json
[
  {
    "id": "da-alert-001",
    "incidentId": "INC-2026-0318-001",
    "title": "Suspicious PowerShell command execution",
    "description": "A process was observed executing a suspicious PowerShell command with encoded content that may indicate malicious activity.",
    "severity": "Critical",
    "status": "New",
    "category": "Execution",
    "classification": null,
    "determination": null,
    "assignedTo": null,
    "creationTime": "2026-03-18T08:00:00.000Z",
    "lastUpdateTime": "2026-03-18T08:05:00.000Z",
    "resolvedTime": null,
    "detectionSource": "EDR",
    "threatFamilyName": "Cobalt Strike",
    "mitreTechniques": ["T1059.001", "T1027"],
    "entities": [
      {
        "entityType": "Process",
        "processId": 4832,
        "processCommandLine": "powershell.exe -enc SQBFAFgAIAAoAE4AZQB3AC0ATwBiAGoAZQBjAHQA...",
        "fileName": "powershell.exe",
        "sha256": "abc123def456abc123def456abc123def456abc123def456abc123def456abcd"
      },
      {
        "entityType": "Machine",
        "deviceId": "dev-001",
        "hostname": "WKS-FINANCE01",
        "osPlatform": "Windows10",
        "riskScore": "High"
      }
    ]
  },
  {
    "id": "da-alert-002",
    "incidentId": "INC-2026-0318-001",
    "title": "Lateral movement via WMI detected",
    "description": "WMI was used to execute commands on a remote machine, which may indicate lateral movement.",
    "severity": "High",
    "status": "InProgress",
    "category": "LateralMovement",
    "classification": "TruePositive",
    "determination": "Malware",
    "assignedTo": "analyst@contoso.com",
    "creationTime": "2026-03-18T08:10:00.000Z",
    "lastUpdateTime": "2026-03-18T09:00:00.000Z",
    "resolvedTime": null,
    "detectionSource": "EDR",
    "threatFamilyName": null,
    "mitreTechniques": ["T1047"],
    "entities": [
      {
        "entityType": "Process",
        "processId": 2108,
        "processCommandLine": "wmic /node:10.0.1.55 process call create \"cmd.exe /c whoami\"",
        "fileName": "WMIC.exe",
        "sha256": "def456abc123def456abc123def456abc123def456abc123def456abc123defg"
      }
    ]
  },
  {
    "id": "da-alert-003",
    "incidentId": "INC-2026-0318-002",
    "title": "Credential dumping tool detected",
    "description": "A known credential dumping tool (mimikatz) was detected running in memory.",
    "severity": "High",
    "status": "New",
    "category": "CredentialAccess",
    "classification": null,
    "determination": null,
    "assignedTo": null,
    "creationTime": "2026-03-18T09:30:00.000Z",
    "lastUpdateTime": "2026-03-18T09:30:00.000Z",
    "resolvedTime": null,
    "detectionSource": "EDR",
    "threatFamilyName": "Mimikatz",
    "mitreTechniques": ["T1003.001", "T1003.002"],
    "entities": [
      {
        "entityType": "Process",
        "processId": 6712,
        "processCommandLine": "sekurlsa::logonpasswords",
        "fileName": "unknown.exe",
        "sha256": "789abc123def456abc123def456abc123def456abc123def456abc123def4567"
      },
      {
        "entityType": "Machine",
        "deviceId": "dev-003",
        "hostname": "DC01",
        "osPlatform": "WindowsServer2022",
        "riskScore": "High"
      }
    ]
  },
  {
    "id": "da-alert-004",
    "incidentId": "INC-2026-0318-003",
    "title": "Suspicious network connection to known C2 server",
    "description": "An outbound connection was detected to a known command-and-control infrastructure IP address.",
    "severity": "Medium",
    "status": "New",
    "category": "CommandAndControl",
    "classification": null,
    "determination": null,
    "assignedTo": null,
    "creationTime": "2026-03-18T10:00:00.000Z",
    "lastUpdateTime": "2026-03-18T10:00:00.000Z",
    "resolvedTime": null,
    "detectionSource": "EDR",
    "threatFamilyName": null,
    "mitreTechniques": ["T1071.001"],
    "entities": [
      {
        "entityType": "Ip",
        "address": "185.220.101.34",
        "location": {
          "countryCode": "DE",
          "city": "Frankfurt"
        }
      }
    ]
  },
  {
    "id": "da-alert-005",
    "incidentId": "INC-2026-0318-004",
    "title": "Anomalous sign-in from unfamiliar location",
    "description": "A user signed in from an IP address in a country not previously seen for this account.",
    "severity": "Medium",
    "status": "Resolved",
    "category": "InitialAccess",
    "classification": "FalsePositive",
    "determination": "NotMalicious",
    "assignedTo": "analyst@contoso.com",
    "creationTime": "2026-03-18T06:00:00.000Z",
    "lastUpdateTime": "2026-03-18T07:30:00.000Z",
    "resolvedTime": "2026-03-18T07:30:00.000Z",
    "detectionSource": "AADIdentityProtection",
    "threatFamilyName": null,
    "mitreTechniques": [],
    "entities": [
      {
        "entityType": "User",
        "accountName": "jsmith@contoso.com",
        "sid": "S-1-5-21-1234567890-123456789-1234567890-1001"
      }
    ]
  },
  {
    "id": "da-alert-006",
    "incidentId": null,
    "title": "Potential data exfiltration via DNS tunneling",
    "description": "Unusually high volume of DNS queries to a single domain detected, which may indicate DNS tunneling.",
    "severity": "Medium",
    "status": "New",
    "category": "Exfiltration",
    "classification": null,
    "determination": null,
    "assignedTo": null,
    "creationTime": "2026-03-18T11:00:00.000Z",
    "lastUpdateTime": "2026-03-18T11:00:00.000Z",
    "resolvedTime": null,
    "detectionSource": "CustomDetection",
    "threatFamilyName": null,
    "mitreTechniques": ["T1048.003"],
    "entities": [
      {
        "entityType": "Machine",
        "deviceId": "dev-005",
        "hostname": "WKS-DEV03",
        "osPlatform": "Windows11",
        "riskScore": "Medium"
      }
    ]
  },
  {
    "id": "da-alert-007",
    "incidentId": null,
    "title": "Suspicious scheduled task creation",
    "description": "A scheduled task was created that runs a script from an unusual location.",
    "severity": "Low",
    "status": "New",
    "category": "Persistence",
    "classification": null,
    "determination": null,
    "assignedTo": null,
    "creationTime": "2026-03-18T12:00:00.000Z",
    "lastUpdateTime": "2026-03-18T12:00:00.000Z",
    "resolvedTime": null,
    "detectionSource": "EDR",
    "threatFamilyName": null,
    "mitreTechniques": ["T1053.005"],
    "entities": [
      {
        "entityType": "Process",
        "processId": 1024,
        "processCommandLine": "schtasks /create /tn UpdateCheck /tr C:\\ProgramData\\update.bat /sc daily",
        "fileName": "schtasks.exe",
        "sha256": null
      }
    ]
  },
  {
    "id": "da-alert-008",
    "incidentId": null,
    "title": "Antivirus exclusion added via command line",
    "description": "A process added an antivirus exclusion path via command line which may indicate defense evasion.",
    "severity": "Low",
    "status": "Resolved",
    "category": "DefenseEvasion",
    "classification": "TruePositive",
    "determination": "SecurityTesting",
    "assignedTo": "admin@contoso.com",
    "creationTime": "2026-03-18T13:00:00.000Z",
    "lastUpdateTime": "2026-03-18T14:00:00.000Z",
    "resolvedTime": "2026-03-18T14:00:00.000Z",
    "detectionSource": "EDR",
    "threatFamilyName": null,
    "mitreTechniques": ["T1562.001"],
    "entities": [
      {
        "entityType": "Process",
        "processId": 3456,
        "processCommandLine": "Set-MpPreference -ExclusionPath 'C:\\Tools'",
        "fileName": "powershell.exe",
        "sha256": null
      }
    ]
  },
  {
    "id": "da-alert-009",
    "incidentId": null,
    "title": "PUA detected: Potentially unwanted application",
    "description": "A potentially unwanted application was detected by Windows Defender Antivirus.",
    "severity": "Informational",
    "status": "Resolved",
    "category": "Malware",
    "classification": "TruePositive",
    "determination": "UnwantedSoftware",
    "assignedTo": null,
    "creationTime": "2026-03-18T14:30:00.000Z",
    "lastUpdateTime": "2026-03-18T14:35:00.000Z",
    "resolvedTime": "2026-03-18T14:35:00.000Z",
    "detectionSource": "WindowsDefenderAv",
    "threatFamilyName": "PUA:Win32/ICBrowserAssistant",
    "mitreTechniques": [],
    "entities": [
      {
        "entityType": "File",
        "fileName": "assistant_setup.exe",
        "filePath": "C:\\Users\\jdoe\\Downloads\\assistant_setup.exe",
        "sha256": "aaa111bbb222ccc333ddd444eee555fff666aaa111bbb222ccc333ddd444eee5"
      }
    ]
  },
  {
    "id": "da-alert-010",
    "incidentId": null,
    "title": "Windows Defender real-time protection disabled",
    "description": "Real-time protection was disabled on a device. This is informational and may be due to a policy change.",
    "severity": "Informational",
    "status": "New",
    "category": "General",
    "classification": null,
    "determination": null,
    "assignedTo": null,
    "creationTime": "2026-03-18T15:00:00.000Z",
    "lastUpdateTime": "2026-03-18T15:00:00.000Z",
    "resolvedTime": null,
    "detectionSource": "WindowsDefenderAv",
    "threatFamilyName": null,
    "mitreTechniques": [],
    "entities": [
      {
        "entityType": "Machine",
        "deviceId": "dev-010",
        "hostname": "WKS-TEST01",
        "osPlatform": "Windows11",
        "riskScore": "Low"
      }
    ]
  }
]
